# etc/cron.d/clean-system
# Clean old System logs
15 04 * * * root /usr/local/bin/clean-system >/dev/null 2>&1
