# etc/cron.d/fun-stuff
#Lets have fortune and cowsay create an motd file every hour
15 * * * * root /usr/local/bin/update-motd >> /etc/motd
