# etc/cron.d/update-rasbian
#Update config from github
15 12 * * * root usr/local/bin/update-rasbian >/dev/null 2>&1
