# etc/cron.d/update-raspbian
#Update configs from github and update the OS at nightly 12:15
15 00 * * * root /usr/local/bin/update-raspbian >/dev/null 2>&1
