# etc/profile.d/raspbian.sh
#!/usr/bin/env bash
alias updateconfigs="sudo update=yes bash -c $(curl -LSs https://github.com/casjay-base/raspbian/raw/master/install.sh)"
