# etc/profile.d/raspbian.sh
#!/usr/bin/env bash
alias updateconfigs="sudo update=yes bash -c "$(curl -LsS https://github.com/casjay/raspbian/raw/master/install.sh)""
